# crontab
0 6 * * * cd /app && python main.py sync >> /var/log/whoop-sync.log 2>&1
